# systemd/system/deck-dictate-l4.service
[Unit]
Description=L4 hold-to-talk dictation bridge
After=systemd-user-sessions.service

[Service]
Type=simple
Environment=SUDO_USER=deck
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/opt/deck-dictate/hold_l4_warm.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
